# systemd/tacacs.target
[Unit]
Description=target to start/stop all tacacs servers at once

[Install]
WantedBy=multi-user.target
